# deploy/provenance.service
[Unit]
Description=Provenance apphost stack (docker compose)
ConditionPathExists=/etc/is-instapphost
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
WorkingDirectory=%h/provenance/deploy
ExecStartPre=-/bin/rm -f /srv/appsockets/provenance/main/app.sock
ExecStartPre=-/usr/bin/docker compose -f compose.apphost.yaml down
# --scale worker=N runs N ingest worker processes so ingest CPU can use more
# than one core (each Node process is single-threaded). Tune N to the box's
# fair share; each worker is capped at mem_limit in the compose file.
ExecStart=/usr/bin/docker compose -f compose.apphost.yaml up --scale worker=2
ExecStop=/usr/bin/docker compose -f compose.apphost.yaml down

[Install]
WantedBy=default.target
